# CTF-2023/40-re-zero/env/src/supervisord.conf
; Documentation of this file format -> http://supervisord.org/configuration.html                                                                                                                                  

[supervisord]
logfile=/dev/stdout
logfile_maxbytes=0
nodaemon=true
loglevel=info
pidfile=/var/run/supervisor/supervisord.pid


[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"
stopasgroup=true
autostart=true
autorestart=true
startsecs=0
startretries=0

;Logs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0

[program:backend]
directory=/app
command=/usr/local/bin/php start.php start
stopasgroup=true
autostart=true
autorestart=true
startsecs=0
startretries=0
user=nobody

;Logs
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stdout_logfile_maxbytes=0